# ansible/roles/edge/files/kasm-api.service
[Unit]
Description=KASM Container Management API
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/utgard-api
ExecStart=/usr/bin/python3 /opt/utgard-api/kasm-api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
