# ansible/roles/openrelik/templates/openrelik-yara-rules.timer.j2
[Unit]
Description=Schedule OpenRelik YARA rules updates

[Timer]
OnCalendar={{ openrelik_yara_rules_update_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
